# modules/nvidia-suspend/files/etc/systemd/system/systemd-suspend.service.d/nvidia-override.conf
[Service]
# Do not freeze user sessions during suspend. Without this, systemd waits up to
# 90 seconds trying to freeze GNOME's session while it holds the NVIDIA driver,
# then suspend fails or takes forever.
Environment=SYSTEMD_SLEEP_FREEZE_USER_SESSIONS=false
