# modules/ssh-sleep-inhibitor/files/etc/systemd/system/ssh-sleep-inhibitor.service
[Unit]
Description=Block sleep while SSH sessions are active

[Service]
Type=simple
ExecStart=/usr/bin/systemd-inhibit \
    --what=sleep \
    --who=SSH \
    --why="SSH session is active" \
    --mode=block \
    /usr/bin/sleep infinity
